# log-tailer.service
[Unit]
Description=Minecraft Log Tailer (Discord notifications)
After=network.target wings.service

[Service]
Type=simple
ExecStart=/usr/bin/python3 /opt/mc-tools/log-tailer.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
